# my-machine-backup.service
[Unit]
Description=My Machine Backup

[Service]
ExecStart=/home/user/bin/backitup.sh
